# ha-addon/rootfs/etc/s6-overlay/s6-rc.d/tuya-dashboard/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Tuya Dashboard..."

cd /app

export NODE_ENV=production
export HOSTNAME=0.0.0.0
export PORT=3000

exec node server.js
